# supervisord.conf

[unix_http_server]
file=/tmp/supervisor.sock

[supervisord]
nodaemon=true
logfile=/var/log/supervisord/supervisord.log 
pidfile=/var/run/supervisord.pid

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[program:mongodb]
autorestart=true
command=/usr/bin/mongod
stderr_logfile = /var/log/supervisord/mongodb.log
stdout_logfile = /var/log/supervisord/mongodb.log

[program:my-slate-api]
autorestart=true
command=node 
stderr_logfile = /var/log/supervisord/my-slate-api.log
stdout_logfile = /var/log/supervisord/my-slate-api.log

[program:my-slate]
autorestart=true
command=/usr/bin/mongod
stderr_logfile = /var/log/supervisord/my-slate.log
stdout_logfile = /var/log/supervisord/my-slate.log